# webapp/scripts/systemd/docklite-agent.service
[Unit]
Description=DockLite Agent
After=network.target docker.service
Requires=docker.service

[Service]
Type=simple
User=docklite
Group=docklite
WorkingDirectory=__INSTALL_DIR__
EnvironmentFile=/etc/docklite/docklite-agent.env
ExecStart=__INSTALL_DIR__/bin/docklite-agent
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
